# rivnet.cron
@reboot root /sbin/iptables-restore /etc/iptables/rules.iptables
